# picofoxweb.service
[Unit]
Description=PICOFoxweb Lightweight Web Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/PICOFoxweb 8000
WorkingDirectory=/var/www/foxweb-jail
RootDirectory=/var/www/foxweb-jail
User=root
Group=root
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target